# packaging/app-base.cron
*/5 * * * * root LANG=en_US /usr/clearos/apps/base/deploy/servicewatch >/dev/null
